# etc/systemd/system/vncserver@.service.d/10-wait-for-home.conf
[Unit]
After=network-online.target sssd.service autofs.service remote-fs.target
Wants=network-online.target autofs.service
