# etc/systemd/system/bootc-update.timer
[Unit]
Description=Weekly bootc image upgrade timer

[Timer]
OnCalendar=Sun *-*-* 04:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
